# deploy/asean-hub-backup.service
[Unit]
Description=Asean English Hub encrypted backup
Wants=asean-english-hub.service

[Service]
Type=oneshot
EnvironmentFile=-/docker/asean-english-hub/.env
WorkingDirectory=/docker/asean-english-hub
ExecStart=/bin/sh /docker/asean-english-hub/scripts/backup-encrypted.sh
User=root
Group=root
