# deploy/asean-hub-backup.timer
[Unit]
Description=Run Asean English Hub encrypted backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
